# USB-AUTO-BACKUP/60-usb-autocopy.rules
ACTION=="add", SUBSYSTEM=="usb", ENV{ID_BUS}=="usb", RUN+="/home/sam/programming/scripts/USB-AUTO-BACKUP/usb-auto-backup.py %k"